# python/systemd/ai.service
[Unit]
Description=AI service
PartOf=evolutek.target

[Service]
ExecStart=ai
CPUSchedulingPolicy=rr

[Install]
WantedBy=multi-user.target
