# python/systemd/basket.service
[Unit]
Description=Basket service
PartOf=evolutek.target

[Service]
ExecStart=basket
CPUSchedulingPolicy=rr

[Install]
WantedBy=multi-user.target
