# python/systemd/buzzer.service
[Unit]
Description=Buzzer service

[Service]
Type=simple
ExecStart=/usr/bin/buzzer

[Install]
WantedBy=multi-user.target

