# python/systemd/cs_echo.service
[Unit]
Description=cs_echo service

[Service]
Type=simple
ExecStart=/usr/bin/cs_echo

[Install]
WantedBy=multi-user.target

