# python/systemd/homologation_pal.service
[Unit]
Description=homologation_pal service

[Service]
Type=simple
ExecStart=/usr/bin/homologation_pal

[Install]
WantedBy=multi-user.target

