# python/systemd/ia_pmi.service
[Unit]
Description=ia_pmi service

[Service]
Type=simple
ExecStart=/usr/bin/ia_pmi

[Install]
WantedBy=multi-user.target

