# python/systemd/leds.service
[Unit]
Description=leds service

[Service]
Type=simple
ExecStart=/usr/bin/leds

[Install]
WantedBy=multi-user.target

