# python/systemd/match.service
[Unit]
Description=Match service
PartOf=evolutek.target

[Service]
ExecStart=match
CPUSchedulingPolicy=rr

[Install]
WantedBy=multi-user.target
