# python/systemd/restart.service
[Unit]
Description=Restart service
PartOf=evolutek.target

[Service]
ExecStart=restart
CPUSchedulingPolicy=rr

[Install]
WantedBy=multi-user.target
