# python/systemd/robots_monitor.service
[Unit]
Description=robot_monitor service

[Service]
Type=simple
ExecStart=/usr/bin/robots_monitor

[Install]
WantedBy=multi-user.target

