# python/systemd/tracking.service
[Unit]
Description=tracking service

[Service]
Type=simple
ExecStart=/usr/bin/tracking

[Install]
WantedBy=multi-user.target

