# python/systemd/trajman.service
[Unit]
Description=Trajman service
PartOf=evolutek.target

[Service]
ExecStart=trajman
CPUSchedulingPolicy=rr

[Install]
WantedBy=multi-user.target
