# ax/systemd/ax.service
[Unit]
Description=AX12 Service

[Service]
ExecStart=/usr/bin/python3 $PREFIX/usr/lib/cs_services/ax/ax.py

[Install]
WantedBy=evolutek.target
